# hack/docker/s6/okmall/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=okmall
exec gosu $USER /app/okmall/okmall serve --config /data/okmall/conf/app.toml
